# udev/95-monitor-hotplug.rules
#Rule for executing commands when an external screen is plugged in.
KERNEL=="card0", ACTION=="change", RUN+="/home/tom/Scripts/monitor-hotplug.sh"
